# watcher.service
[Unit]
Description=Watcher service

[Service]
Type=simple
ExecStart=/home/numeric/watch_parse_send_data/watcher

[Install]
WantedBy=multi-user.target
